# infra/assets/services/openclaw-backup.service
[Unit]
Description=OpenClaw backup job

[Service]
Type=oneshot
ExecStart=%h/bin/openclaw-backup.sh
