# infra/assets/services/workspace-sync.service
[Unit]
Description=Workspace git sync job

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'cd ~/.openclaw/workspace-main && git add -A && git diff --cached --quiet || git commit -m "auto-sync $(date -u +%%Y-%%m-%%dT%%H:%%M:%%SZ)" && git push'
