# infra/assets/services/workspace-sync.timer
[Unit]
Description=Daily workspace git sync

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
